# jekyll-serve.supervisor.conf
[program:jekyll-serve]
autostart=true
autorestart=true
redirect_stderr=true
command=sudo su vagrant -l -c 'bash /vagrant/start-jekyll.sh'
